# infra/systemd/membra-api.service
[Unit]
Description=MEMBRA QR Gateway — FastAPI service
Documentation=https://github.com/membra/membra-qr-gateway
After=network.target
Wants=network.target

[Service]
Type=exec
User=membra
Group=membra
WorkingDirectory=/opt/membra-qr-gateway

# Load secrets from host-level env file (chmod 600, owned by membra)
EnvironmentFile=/etc/membra/app.env

# Override specific vars (non-secret)
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PORT=7860

ExecStart=/opt/membra-qr-gateway/.venv/bin/python -m uvicorn api.main:app \
    --host 0.0.0.0 \
    --port 7860 \
    --workers 4 \
    --log-level info \
    --access-log \
    --proxy-headers \
    --forwarded-allow-ips "*"

ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

# Graceful shutdown timeout
TimeoutStopSec=30s
KillMode=mixed
KillSignal=SIGTERM

# Output to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=membra-api

# ── Security sandboxing ───────────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/membra-qr-gateway/data /opt/membra-qr-gateway/storage
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0027

# ── Resource limits ───────────────────────────────────────────────────────────
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=150%

[Install]
WantedBy=multi-user.target
